# data/xdg-desktop-portal-gtk.service.in
[Unit]
Description=Portal service (GTK/GNOME implementation)

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.gtk
ExecStart=@libexecdir@/xdg-desktop-portal-gtk
